# server/crontab
* * * * * root cd /app && /bin/bash pubrunner.sh >> /var/log/cron.log 2>&1
