# server/crontab
* * * * * root R pubmed_download.R >> /var/log/cron.log 2>&1
